# etc/systemd/system/memcached.service.d/socket.conf
# Configure memcached.service to use socket instead of localhost.

[Service]
ExecStart=
ExecStart=/usr/bin/memcached -a 666 -s /var/run/memcached.sock
