# etc/systemd/system/memcached.service.d/socket.conf
# Configure memcached.service to use socket instead of localhost.

[Service]
User=
User=root
ExecStart=
ExecStart=/usr/bin/memcached -u memcached -a 666 -s /var/run/memcached/memcached.sock
